# registry/etc/systemd/system/docker.service.d/execstart.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --data-root /home/kayw/fshare/.dockercache
