# registry/etc/local.d/vboxsf.stop
#!/bin/sh -e
umount /home/kayw/share
umount /media/sf_shared
exit 0
